# rootfs/etc/inittab
::sysinit:/bin/mount -t proc proc /proc
::respawn:-/bin/sh
